# sonar/supervisor/cron.conf
[program:cron]
command=/usr/sbin/crond -f -L 0 -d 0 -c /etc/crontabs/ -L /var/log/supervisor/crond
autorestart=true
