# sonar/supervisor/sonarqube.conf
[program:sonarqube]
command=/opt/sonarqube/bin/run.sh
autorestart=true
stdout_logfile=/var/log/supervisor/sonarqube
stderr_logfile=/var/log/supervisor/sonarqube
