# sonar/crontab
0       0       *       *       *       /opt/sonarqube/bin/backup.sh
