# sonar/crontab
0       0       *       *       *       /opt/sonarqube/bin/sonarqube_backup
