# postgresql/supervisor/postgresql.conf
[program:postgresql]
command=/docker-entrypoint.sh
autorestart=true
stdout_logfile=/var/log/supervisor/postgresql
stderr_logfile=/var/log/supervisor/postgresql
